# scripts/com.boathire24.blogger-drain.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.boathire24.blogger-drain</string>

    <!-- Drains the Blogger backlog (posts + pages) a safe batch at a time.
         The publisher aborts cleanly when Blogger write-blocks (403/quota) and
         resumes next run, so this just keeps topping up whenever quota allows. -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/master/boat-rental-platform/scripts/post_blogger.py</string>
        <string>--limit</string>
        <string>8</string>
    </array>

    <!-- Every 4 hours -->
    <key>StartInterval</key>
    <integer>14400</integer>

    <key>WorkingDirectory</key>
    <string>/Users/master/boat-rental-platform</string>

    <key>StandardOutPath</key>
    <string>/Users/master/boat-rental-platform/logs/blogger_drain.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/master/boat-rental-platform/logs/blogger_drain.stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
